# pi/etc/register.cron
@hourly /home/pi/garage-control/pi/scripts/register.sh > /home/pi/garage-control/pi/logs/register.log 2>&1
